# linux/r2m-server.socket
[Unit]
Description=r2m-server socket

[Socket]
FileDescriptorName=r2m-server
ListenStream=80
NoDelay=true